# diy/cosmetics-diy.service
[Unit]
Description=Cosmetics DIY Web Application
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/app  # 替换为实际的应用路径
Environment=PATH=/path/to/your/venv/bin  # 替换为虚拟环境路径
ExecStart=/path/to/your/venv/bin/gunicorn --config gunicorn.conf.py wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

